# LaunchAgents/com.alexberry.theme-switcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexberry.theme-switcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/dark-notify</string>
        <string>-c</string>
        <string>/Users/alex.berry/Projects/highwaybobbery/dotfiles/master/scripts/theme-switcher</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/alex.berry/Library/Logs/theme-switcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alex.berry/Library/Logs/theme-switcher.error.log</string>
</dict>
</plist>